# deploy/systemd/lingxi-cleaner.service
[Unit]
Description=LingXi expired file cleanup

[Service]
Type=oneshot
User=lingxi
WorkingDirectory=/opt/lingxi/apps/api
EnvironmentFile=/opt/lingxi/apps/api/.env
ExecStart=/opt/lingxi/apps/api/.venv/bin/python -m app.tasks.cleanup
